# ADAS_Virtual_ECU/sw_app/platform/sdk/rtos/vision_apps/apps/edt_server/edtserver.service
[Unit]
Description="ED Test server"
After=systemd-modules-load.service

[Service]
Restart=always
RestartSec=1
ExecStart=/opt/vision_apps/edt_server.out &

[Install]
WantedBy=multi-user.target
